# 003-Docker-Cron-n-Apache/supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
command=/usr/sbin/crond

[program:httpd]
command=/usr/sbin/httpd -k start
